# sfml-platform.service
[Unit]
Description=SFML Platform (compat wrapper for SHML Platform startup)
Documentation=https://github.com/axelofwar/shml-platform
After=docker.service network-online.target tailscaled.service
Requires=docker.service
Wants=network-online.target tailscaled.service
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=oneshot
RemainAfterExit=yes
User=axelofwar
Group=docker
WorkingDirectory=/opt/shml-platform
ExecStartPre=/bin/sleep 10
ExecStartPre=/bin/bash -c 'until docker info >/dev/null 2>&1; do sleep 2; done'
ExecStartPre=/bin/bash -c 'until tailscale status >/dev/null 2>&1; do sleep 2; done || true'
ExecStart=/opt/shml-platform/start_all_safe.sh start
ExecStop=/opt/shml-platform/stop_all.sh
Restart=on-failure
RestartSec=30
TimeoutStartSec=600
TimeoutStopSec=120
Environment=HOME=/home/$USER
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=OPTIONAL_AI_STACK_ENABLED=false
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shml-platform

[Install]
WantedBy=multi-user.target
