# deploy/systemd/shml-platform.service
[Unit]
Description=SFML Platform (MLflow + Ray + Auth Services)
Documentation=https://github.com/axelofwar/shml-platform
After=docker.service network-online.target tailscaled.service shml-tailscale-funnel.service
Requires=docker.service shml-tailscale-funnel.service
Wants=network-online.target tailscaled.service

# Give time for network and Tailscale to fully initialize
# This is critical for OAuth/FusionAuth which need Tailscale Funnel URLs
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=oneshot
RemainAfterExit=yes
User=${PLATFORM_USER:-axelofwar}
Group=docker

# Working directory for docker-compose
WorkingDirectory=/opt/shml-platform

# Wait for dependencies before starting
ExecStartPre=/bin/sleep 10
ExecStartPre=/bin/bash -c 'until docker info >/dev/null 2>&1; do sleep 2; done'
ExecStartPre=/bin/bash -c 'until tailscale status >/dev/null 2>&1; do sleep 2; done'

# Start all services using the safe startup script
# The script handles all internal dependencies (Postgres -> FusionAuth -> OAuth2 Proxy -> Services)
ExecStart=/opt/shml-platform/start_all_safe.sh start

# Stop all services gracefully
ExecStop=/opt/shml-platform/stop_all.sh

# Restart policy
Restart=on-failure
RestartSec=30

# Timeout for long-running startup (FusionAuth can take 2-3 minutes)
TimeoutStartSec=600
TimeoutStopSec=120

# Environment
Environment=HOME=/home/$USER
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=OPTIONAL_AI_STACK_ENABLED=false

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shml-platform

[Install]
WantedBy=multi-user.target
