# inference/app/native_training/systemd/shml-training.service
[Unit]
Description=SHML Native Training Job
Documentation=file:///home/axelofwar/Projects/shml-platform/inference/app/native_training/README.md
After=network-online.target docker.service
Wants=network-online.target
# Can be started independently or by coordinator
PartOf=shml-training-coordinator.service

[Service]
Type=simple
User=axelofwar
Group=axelofwar

# Working directory
WorkingDirectory=/home/axelofwar/Projects/shml-platform/inference/app/native_training

# Environment for GPU training
Environment=PYTHONPATH=/home/axelofwar/Projects/shml-platform/inference/app:/home/axelofwar/Projects/shml-platform/inference/training_library
Environment=CUDA_VISIBLE_DEVICES=0
Environment=CUDA_DEVICE_ORDER=PCI_BUS_ID
Environment=PYTORCH_CUDA_ALLOC_CONF=max_split_size_mb:512
Environment=MLFLOW_TRACKING_URI=http://172.30.0.11:5000
# CRITICAL: Bypass MPS to allow direct CUDA access alongside MPS-based inference
Environment=CUDA_MPS_PIPE_DIRECTORY=

# Main process - sandboxed training
ExecStart=/home/axelofwar/Projects/shml-platform/inference/app/native_training/sandbox_training.sh \
    --model yolov8n.pt \
    --data wider_face \
    --epochs 100 \
    --batch 16 \
    --mlflow-uri http://172.30.0.11:5000 \
    --experiment SOTA-YOLO-Training

# Graceful shutdown with checkpoint
ExecStop=/bin/kill -TERM $MAINPID
TimeoutStopSec=120

# Restart on failure (but not on success - training is expected to finish)
Restart=on-failure
RestartSec=30

# RESOURCE LIMITS
# -----------------------------------------------------------------------------
# Memory: 58GB max (leaves 6GB for system + inference)
MemoryMax=58G
MemoryHigh=50G

# CPU: Allow 8 cores but lower priority than inference
CPUQuota=800%
CPUWeight=50

# IO: Lower priority for disk operations
IOWeight=50

# Standard output and error
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shml-training

# Security - minimal because sandbox_training.sh handles isolation
NoNewPrivileges=true

# GPU access requires some /dev access
DeviceAllow=char-nvidia rw
DeviceAllow=char-drm rw

[Install]
WantedBy=multi-user.target
