# scripts/monitoring/systemd/autonomous-guard-remediation.timer
[Unit]
Description=Run SHML Autonomous Guard Remediation every 5 minutes

[Timer]
OnBootSec=90s
OnUnitActiveSec=5m
AccuracySec=30s
Persistent=true
Unit=autonomous-guard-remediation.service

[Install]
WantedBy=timers.target
