# scripts/monitoring/systemd/shl-gitlab-sync.service
[Unit]
Description=SHL GitLab board sync — sync T8 pipeline state to GitLab Issues
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/Projects/shml-platform
ExecStart=/usr/bin/env bash -lc '\
    exec bash "$HOME/Projects/shml-platform/scripts/data/update_gitlab_board.sh"'
TimeoutStartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
