# scripts/monitoring/systemd/shl-nano-kanban.timer
[Unit]
Description=SHL Nano Kanban updater — every 10 minutes

[Timer]
OnBootSec=2m
OnUnitActiveSec=10m
AccuracySec=60s
Persistent=false
Unit=shl-nano-kanban.service

[Install]
WantedBy=timers.target
