# scripts/tailscale-funnel.service
[Unit]
Description=Tailscale Funnel for SFML Platform
Documentation=https://tailscale.com/kb/1223/funnel/
After=network-online.target tailscaled.service
Wants=network-online.target
Requires=tailscaled.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/tailscale funnel --bg --https=443 http://127.0.0.1:80
ExecStop=/usr/bin/tailscale funnel --https=443 off
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
